# deploy/linux/easyfire-bookkeeping-guardian.timer
[Unit]
Description=Run EasyFire Bookkeeping Guardian every 30 seconds

[Timer]
OnBootSec=120s
OnUnitActiveSec=30s
AccuracySec=2s
RandomizedDelaySec=3s
Persistent=false
Unit=easyfire-bookkeeping-guardian.service

[Install]
WantedBy=timers.target
