# deploy/linux/easyfire-bookkeeping-stack.service
[Unit]
Description=EasyFire Bookkeeping existing-container startup gate
Documentation=file:/opt/easyfire-bookkeeping/current/docs/easyfire/LINUX_VM_RUNBOOK.md
Requires=docker.service
PartOf=docker.service
After=docker.service network-online.target
Wants=network-online.target
ConditionPathExists=/etc/easyfire-bookkeeping/production.env
ConditionPathExists=/etc/easyfire-bookkeeping/deployment-plan.json
ConditionPathExists=/etc/easyfire-bookkeeping/checkpoint-manifest.json
ConditionPathExists=/etc/easyfire-bookkeeping/migration-receipt.json
ConditionPathExists=/etc/easyfire-bookkeeping/deployment-receipt.json
ConditionPathExists=/etc/easyfire-bookkeeping/runtime-manifest.json
ConditionPathExists=/etc/easyfire-bookkeeping/runtime-identity-evidence.json
ConditionPathExists=!/etc/easyfire-bookkeeping/rollback.lock
ConditionPathExists=/opt/easyfire-bookkeeping/current/docker-compose.prod.yml
ConditionPathExists=/opt/easyfire-bookkeeping/current/deploy/linux/docker-compose.vm.yml

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=600
UMask=0077
RuntimeDirectory=easyfire-bookkeeping-stack
RuntimeDirectoryMode=0700
ExecStartPre=/usr/local/bin/node /opt/easyfire-bookkeeping/current/scripts/production/linux-deploy-candidate.mjs --verify-existing --plan /etc/easyfire-bookkeeping/deployment-plan.json
ExecStart=/usr/bin/docker compose --project-name easyfire-bookkeeping-prod --file /opt/easyfire-bookkeeping/current/docker-compose.prod.yml --file /opt/easyfire-bookkeeping/current/deploy/linux/docker-compose.vm.yml --file /opt/easyfire-bookkeeping/current/deploy/linux/docker-compose.candidate.yml --env-file /etc/easyfire-bookkeeping/production.env start --wait --wait-timeout 180 mysql redis
ExecStart=/usr/bin/docker compose --project-name easyfire-bookkeeping-prod --file /opt/easyfire-bookkeeping/current/docker-compose.prod.yml --file /opt/easyfire-bookkeeping/current/deploy/linux/docker-compose.vm.yml --file /opt/easyfire-bookkeeping/current/deploy/linux/docker-compose.candidate.yml --env-file /etc/easyfire-bookkeeping/production.env start --wait --wait-timeout 180 gotenberg server webapp envoy
ExecStartPost=/usr/bin/touch /run/easyfire-bookkeeping-stack/ready

[Install]
WantedBy=multi-user.target
